# scripts/kryos-dryrun.timer
[Unit]
Description=KryOs dry-run timer (A/B test) - runs every 10 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
AccuracySec=1s
Unit=kryos-dryrun.service

[Install]
WantedBy=timers.target
